# deployer/deploy/config/update-motd.d/80-who
#!/bin/sh
echo "Uptime:"
/usr/bin/uptime
echo "Previous logins:"
lastlog -t 30
echo "Currently logged in:"
who